# ticker-controller.service
[Unit]
Description=Sports Ticker Controller (LED Display)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 -m ticker_controller
Restart=always
RestartSec=5
StandardOutput=append:/home/pi/ticker.log
StandardError=append:/home/pi/ticker.log

[Install]
WantedBy=multi-user.target
